# tg-agent.service
[Unit]
Description=Telegram AI Agent → DeepSeek → Google Sheets
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/tg-agent
EnvironmentFile=/home/pi/tg-agent/.env
ExecStart=/home/pi/tg-agent/.venv/bin/python telegram_agent.py
Restart=always
RestartSec=10

# Логи через journald
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
